# server/systemd/gmail-watcher.timer
[Unit]
Description=Run Gmail Watcher every 15 minutes

[Timer]
OnBootSec=5min
OnUnitActiveSec=15min
Unit=gmail-watcher.service

[Install]
WantedBy=timers.target
