# server/systemd/gmail-watcher.service
[Unit]
Description=Gmail Reply Detection Service for Notiapply
After=network.target postgresql.service

[Service]
Type=oneshot
User=%REPLACE_USER%
WorkingDirectory=%REPLACE_WORKDIR%
Environment="DATABASE_URL=%REPLACE_DB_URL%"
ExecStart=%REPLACE_VENV%/bin/python gmail_watcher.py --check --db-url $DATABASE_URL

# Logging
StandardOutput=append:/var/log/notiapply/gmail-watcher.log
StandardError=append:/var/log/notiapply/gmail-watcher-error.log

[Install]
WantedBy=multi-user.target
